# .config/autostart/devmon.desktop
[Desktop Entry]
Name=devmon
Comment=Automount USB devices
Exec=devmon
Terminal=false
Type=Application
StartupNotify=false
