# .config/autostart/mpd.desktop
[Desktop Entry]
Name=mpd
Comment=Music Player Daemon
Exec=mpd
Terminal=false
Type=Application
StartupNotify=false
