# .config/autostart/parcellite.desktop
[Desktop Entry]
Name=Parcellite
Comment=Clipboard management
Exec=parcellite
Terminal=false
Type=Application
StartupNotify=false
