# .config/autostart/playerctld.desktop
[Desktop Entry]
Name=playerctld
Comment=Music
Exec=playerctld daemon
Terminal=false
Type=Application
StartupNotify=false
