# .config/autostart/radeon-profile.desktop
[Desktop Entry]
Name=radeon-profile
Comment=Shut your fans up
Exec=sudo radeon-profile
Terminal=false
Type=Application
StartupNotify=false
