# .config/autostart/redshift.desktop
[Desktop Entry]
Name=Redshift
Comment=Protect your eyes
Exec=redshift
Terminal=false
Type=Application
StartupNotify=false
