# systemd/.config/systemd/user/geoclue-agent.service
[Unit]
Description=I've got a raging (geo)clue

[Service]
ExecStart=/usr/lib/geoclue-2.0/demos/agent

[Install]
WantedBy=default.target
